# gunicorn/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=unix:/run/gunicorn.sock

[Install]
WantedBy=sockets.target